# supervisord.conf
; Process supervisor for the TurboVid container.
;
; Two long-running children, one config:
;   [program:web]    — uvicorn serving the FastAPI app on $BULKVID_PORT.
;   [program:worker] — the queue-draining ``python -m bulkvid.worker`` loop.
;
; Both autorestart so a crash doesn't quietly take the service down. Both
; log straight to the container's stdout/stderr so HuggingFace's "Logs"
; tab (and any other Docker host) sees the application logs without
; volume mounts.
;
; Plan: ``_plans/2026-06-04-migrate-to-hf-spaces-turso.md`` §Architecture.

[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:web]
command=/usr/local/bin/uvicorn bulkvid.main:app --host 0.0.0.0 --port %(ENV_BULKVID_PORT)s
directory=/app
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=15
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

[program:worker]
command=/usr/local/bin/python -m bulkvid.worker
directory=/app
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20
